# timers/backup_cloud.timer
[Unit]
Description=Backs up to the cloud daily

[Timer]
OnCalendar=*-*-* __BACKUP_CLOUD_TIME__
Persistent=true

[Install]
WantedBy=timers.target