# services/check_mount.service
[Unit]
Description=Check USB Mount

[Service]
ExecStart=/usr/local/bin/check_mount.sh

[Install]
WantedBy=multi-user.target